# trading/deploy/aws/systemd/sfo-kalshi-paper-scan.timer
[Unit]
Description=Run SFO Kalshi paper scans every 15 minutes

[Timer]
OnCalendar=*-*-* *:00,15,30,45
Persistent=true
Unit=sfo-kalshi-paper-scan.service

[Install]
WantedBy=timers.target
