# trading/deploy/aws/systemd/sfo-kalshi-paper-monitor.timer
[Unit]
Description=Run SFO Kalshi paper exit monitor every 5 minutes

[Timer]
OnCalendar=*-*-* *:02,07,12,17,22,27,32,37,42,47,52,57
Persistent=true
Unit=sfo-kalshi-paper-monitor.service

[Install]
WantedBy=timers.target
